# src/mpDris2.service.in
[Unit]
Description=mpDris2 - Music Player Daemon D-Bus bridge
After=mpd.service
BindsTo=mpd.service
Conflicts=mpDris2-remote.service

[Service]
Restart=on-failure
ExecStart=@bindir@/mpDris2 --use-journal --no-reconnect
BusName=org.mpris.MediaPlayer2.mpd

[Install]
WantedBy=mpd.service
