# src/mpDris2-remote.service.in
[Unit]
Description=mpDris2 - Music Player Daemon D-Bus bridge
Conflicts=mpDris2.service

[Service]
Restart=on-failure
RestartSec=5
RestartMaxDelaySec=15min
RestartSteps=20
ExecStart=@bindir@/mpDris2 --use-journal --no-reconnect --abort-on-disconnect
BusName=org.mpris.MediaPlayer2.mpd

[Install]
WantedBy=default.target
