# src/mpDris2.service.in
[Unit]
Description=mpDris2 - Music Player Daemon D-Bus bridge

[Service]
Restart=on-failure
ExecStart=@bindir@/mpDris2
BusName=org.mpris.MediaPlayer2.mpd

[Install]
WantedBy=default.target
# WantedBy=daemon.target
